# scripts/.config/systemd/user/wallpaper.service
[Unit]
Description=Set random wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper.sh
Environment=DISPLAY=:0
